# cron/stats-index-cleaner/crontab
* * * * * root env - `cat /root/.env` /src/clean.sh
